# package/dni/control-firmware/files/control_firmware.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=96
boot() {
	is_base=`cat /tmp/orbi_type`
	if [ "x$is_base" = "xBase" ]; then
		echo "control-firmware will be Running when Internet status is UP !!!" > /dev/console
		/usr/sbin/boot-notification > /dev/console &
	fi
}
