# package/dni/base-files_dni/files/etc/init.d/ddns
#!/bin/sh /etc/rc.common

[ $(/bin/config get factory_mode) -eq 1 ] && return

start() {
	/sbin/cmd_ddns start
}

stop() {
	/sbin/cmd_ddns stop
}

restart() {
	/sbin/cmd_ddns restart
}

